# kmonad/kmonad.service
[Unit]
Description=Kmonad (keyboard mapper)

[Service]
ExecStart=/home/umut/dotfiles/kmonad/start-all.sh

[Install]
WantedBy=default.target
